# recipes-core/systemd/files/rootfs-resize.service
[Unit]
Description=One time boot script

[Service]
Type=oneshot
ExecStart=/usr/bin/rootfs-resize.sh

[Install]
WantedBy=multi-user.target
